# ops/ansible/roles/nginx-stack/templates/pressluft-acme-renew.service.j2
[Unit]
Description=Renew ACME certificates for Pressluft
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pressluft
Group=pressluft
Environment=HOME=/var/lib/pressluft
ExecStart=/usr/local/bin/pressluft-acme-renew
