# ops/ansible/roles/nginx-stack/templates/pressluft-acme-renew.timer.j2
[Unit]
Description=Daily ACME renewal for Pressluft

[Timer]
OnBootSec=10m
OnUnitActiveSec=1d
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
